# opt/etc/profile.d/z01_lmod.sh
MODULEPATH=/data/opt/etc/modulefiles:$MODULEPATH

module load anaconda/3.7
